# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:fastapi]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/fastapi_err.log
stdout_logfile=/app/logs/fastapi.log

[program:paper_trading]
command=python scripts/start_paper_trading.py
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/paper_trading_err.log
stdout_logfile=/app/logs/paper_trading.log
